# etc/cron.daily/geoipupdate.sh
#!/bin/bash

if [ -f /usr/bin/geoipupdate ];
	then /usr/bin/geoipupdate
fi
